# CONFIG_FILES/%etc%init%ntp.conf
# ntp daemon
#

description	"NTP daemon"
author		"None"

start on started early-network-time
stop on runlevel [!2345]

expect fork
respawn

exec /usr/sbin/ntpd -p /var/run/ntpd.pid -g -u ntp:ntp
