# CONFIG_FILES/%etc%init%net.conf
# run wired at startup

description	"connect to ethernet"

start on (started dbus and started mountall)
stop on stopping dbus

exec wired
